# etc/supervisor.d/relayer.conf
[program:relayer]
command=run-relayer
autostart=false
autorestart=true
stderr_logfile=var/log/relayer.log
stdout_logfile=var/log/relayer.log
